# deploy/barde-restart.service
[Unit]
Description=Redemarrage quotidien de Barde (hygiene process 24/7)

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart barde
